# src/root/etc/cron.d/check_tools_repo
04 07 * * * root bash /opt/dan_tools/bin/check-repo.sh /opt/dan_tools/ local

