# src/root/etc/cron.d/reboot-reminder
@reboot root echo "Just to let you know that the $(hostname -f) has been rebooted" | mail -s "Server $(hostname) has been rebooted" root
